# root/etc/uci-defaults/30_luci-theme-noodles
#!/bin/sh
uci batch <<-EOF
	set luci.themes.noodlesPurple=/luci-static/noodlePurple
	commit luci
EOF
exit 0
